# ops/moode-overrides/etc/systemd/system/airplay-json-watchdog.service
[Unit]
Description=Watchdog for AirPlay metadata pipeline CPU
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/airplay-json-watchdog.sh
